# logstash/pipeline/main.conf
input {
  beats {
    port => 5044
    host => "0.0.0.0"
  }
}

filter {

  mutate {
    add_field => {
      "received_at" => "%{@timestamp}"
      "soc_version" => "2.0"
    }
  }

  #################################################
  # AUTH LOGS
  #################################################

  if [fields][log_type] == "auth" {

    grok {
      match => {
        "message" => [
          "%{SYSLOGTIMESTAMP:event_time} %{HOSTNAME:reporting_host} %{DATA:program}(?:\[%{POSINT:pid}\])?: %{GREEDYDATA:log_message}",
          "%{TIMESTAMP_ISO8601:event_time} %{HOSTNAME:reporting_host} %{DATA:program}(?:\[%{POSINT:pid}\])?: %{GREEDYDATA:log_message}"
        ]
      }
    }

    #
    # Failed Login
    #

    if "Failed password" in [log_message] {

      mutate {
        add_tag => [
          "failed_login",
          "brute_force_candidate",
          "ATTACK_T1110"
        ]
      }

      grok {
        match => {
          "log_message" => [
            "Failed password for %{DATA:user} from %{IP:src_ip}",
            "Failed password for invalid user %{DATA:user} from %{IP:src_ip}"
          ]
        }
      }
    }

    #
    # Invalid User
    #

    if "Invalid user" in [log_message] {

      mutate {
        add_tag => [
          "invalid_user",
          "recon_candidate",
          "ATTACK_T1110"
        ]
      }

      grok {
        match => {
          "log_message" =>
          "Invalid user %{DATA:user} from %{IP:src_ip}"
        }
      }
    }

    #
    # Successful Login
    #

    if "Accepted password" in [log_message] or
       "Accepted publickey" in [log_message] {

      mutate {
        add_tag => [
          "successful_login",
          "ATTACK_T1078"
        ]
      }

      grok {
        match => {
          "log_message" =>
          "Accepted %{WORD:auth_type} for %{DATA:user} from %{IP:src_ip}"
        }
      }
    }

    #
    # SSH Scanning
    #

    if "Did not receive identification string" in [log_message] {

      mutate {
        add_tag => [
          "ssh_scan",
          "recon_candidate"
        ]
      }
    }

    #
    # Privilege Escalation
    #

    if [log_message] =~ /(?i)sudo/ {

      mutate {
        add_tag => [
          "privilege_event",
          "ATTACK_T1548"
        ]
      }
    }

    #
    # Failed su
    #

    if [log_message] =~ /(?i)FAILED su/ {

      mutate {
        add_tag => [
          "lateral_movement_candidate",
          "ATTACK_T1078"
        ]
      }
    }

    #
    # User Creation
    #

    if [log_message] =~ /(?i)(useradd|adduser)/ {

      mutate {
        add_tag => [
          "user_creation",
          "ATTACK_T1136"
        ]
      }
    }

    #
    # Password Changes
    #

    if [log_message] =~ /(?i)(passwd|chpasswd)/ {

      mutate {
        add_tag => [
          "account_modification",
          "ATTACK_T1098"
        ]
      }
    }

    #
    # Scheduled Tasks
    #

    if [log_message] =~ /(?i)(cron|crontab)/ {

      mutate {
        add_tag => [
          "scheduled_task",
          "ATTACK_T1053"
        ]
      }
    }
  }

  #################################################
  # NGINX ACCESS LOGS
  #################################################

  if [fields][log_type] == "nginx_access" {

    #
    # Scanner Detection
    #

    if [message] =~ /(?i)(sqlmap|nikto|masscan|nmap|wpscan|gobuster|dirbuster|acunetix)/ {

      mutate {
        add_tag => [
          "scanner_detected",
          "web_attack"
        ]
      }
    }

    #
    # SQL Injection
    #

    if [message] =~ /(?i)(union.*select|select.*from|drop\s+table|information_schema|sleep\()/ {

      mutate {
        add_tag => [
          "sql_injection_attempt",
          "web_attack",
          "ATTACK_T1190"
        ]
      }
    }

    #
    # XSS
    #

    if [message] =~ /(?i)(<script>|javascript:|onerror=|onload=)/ {

      mutate {
        add_tag => [
          "xss_attempt",
          "web_attack"
        ]
      }
    }

    #
    # Path Traversal
    #

    if [message] =~ /(\.\.\/)+/ {

      mutate {
        add_tag => [
          "path_traversal_attempt",
          "web_attack"
        ]
      }
    }

    #
    # Sensitive File Access
    #

    if [message] =~ /(?i)(\.env|id_rsa|shadow|passwd|\.git|config\.php)/ {

      mutate {
        add_tag => [
          "sensitive_file_access",
          "web_attack"
        ]
      }
    }

    #
    # Command Injection
    #

    if [message] =~ /(;|\|\||&&)\s*(wget|curl|bash|sh|python)/ {

      mutate {
        add_tag => [
          "command_injection",
          "ATTACK_T1059"
        ]
      }
    }

    #
    # Reverse Shell
    #

    if [message] =~ /(?i)(bash -i|nc -e|python -c.*socket|perl -e.*socket)/ {

      mutate {
        add_tag => [
          "reverse_shell",
          "ATTACK_T1059"
        ]
      }
    }

    #
    # Log4Shell
    #

    if [message] =~ /\$\{jndi:/ {

      mutate {
        add_tag => [
          "log4shell_attempt",
          "ATTACK_T1190"
        ]
      }
    }

    #
    # Web Shell
    #

    if [message] =~ /(?i)(cmd=|exec=|shell=|system\()/ {

      mutate {
        add_tag => [
          "webshell_activity",
          "ATTACK_T1505"
        ]
      }
    }
  }

  #################################################
  # SYSLOG
  #################################################

  if [fields][log_type] == "syslog" {

    if [message] =~ /apparmor="DENIED"/ {

      mutate {
        add_tag => [
          "security_control_event"
        ]
      }
    }

    if [message] =~ /(?i)(wget|curl).*http/ {

      mutate {
        add_tag => [
          "suspicious_download",
          "ATTACK_T1105"
        ]
      }
    }

    if [message] =~ /(?i)(nc -e|bash -i|python -c)/ {

      mutate {
        add_tag => [
          "reverse_shell",
          "ATTACK_T1059"
        ]
      }
    }
  }
}

output {

  file {
    path => "/usr/share/logstash/data/events/%{[fields][log_type]}-%{+YYYY-MM-dd}.log"
    codec => json_lines
  }

  stdout {
    codec => rubydebug
  }
}